# files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}
{%- set service_config = pillar.ffmpeg.server.input[service.name] %}
[program:{{ service.type }}_{{ service.name }}]
directory=/opt/ffmpeg
command=ffserver -f /opt/ffmpeg/{{ service.name }}.conf & {{ service_config.video_format }} -v 2 -r 5 -s {{ service_config.width }}x{{ service_config.height }} -f video4linux2 -i {{ service_config.source }} {{ service_config.bind.host }}:{{ service_config.bind.port }}/{{ service.name }}.{{ service_config.format }}
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
stdout_logfile=/var/log/ffmpeg/{{ service.name }}.log 
stderr_logfile=/var/log/ffmpeg/{{ service.name }}_err.log
user=ffmpeg
